# scripts/com.genesis.factory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Genesis Factory — Auto-start on login.
  
  Install:
    cp scripts/com.genesis.factory.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.genesis.factory.plist
  
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.genesis.factory.plist
    rm ~/Library/LaunchAgents/com.genesis.factory.plist
-->
<dict>
    <key>Label</key>
    <string>com.genesis.factory</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <!-- Update this path to where you cloned genesis-factory -->
        <string>$HOME/genesis-factory/scripts/start.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/genesis-factory.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/genesis-factory.err</string>
</dict>
</plist>
