# scripts/com.genesis.factory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Genesis Factory — Auto-start on login.

  Install (update the path first!):
    FACTORY_DIR="$HOME/genesis-factory"
    sed "s|__FACTORY_DIR__|$FACTORY_DIR|g" scripts/com.genesis.factory.plist > /tmp/gf.plist
    cp /tmp/gf.plist ~/Library/LaunchAgents/com.genesis.factory.plist
    launchctl load ~/Library/LaunchAgents/com.genesis.factory.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.genesis.factory.plist
    rm ~/Library/LaunchAgents/com.genesis.factory.plist
-->
<dict>
    <key>Label</key>
    <string>com.genesis.factory</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <!-- __FACTORY_DIR__ is replaced during install (see instructions above) -->
        <string>__FACTORY_DIR__/scripts/start.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/genesis-factory.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/genesis-factory.err</string>
</dict>
</plist>
